# upstart_wengine.conf
description "Wengine Server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
    /actuator/wengine/src/wengine/wengine
end script
