# upstart_wapour.conf
description "Wapour Server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
    /actuator/wapour/src/wapour/wapour
end script
