# wapour.service
[Unit]
Description=wapour
After=wengine.service

[Service]
Type=simple
ExecStart=/actuator/wapour/src/wapour/wapour
